# scripts/launchd/io.github.auto-model-switcher.router.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.auto-model-switcher.router</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>@INSTALL_DIR@/scripts/run-freerouter.sh</string>
    <string>@INSTALL_DIR@</string>
    <string>@CONFIG_DIR@</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/auto-model-switcher/router.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/auto-model-switcher/router.error.log</string>
</dict>
</plist>
