# scripts/launchd/io.github.auto-model-switcher.update.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>io.github.auto-model-switcher.update</string>
  <key>ProgramArguments</key>
  <array>
    <string>@INSTALL_DIR@/bin/auto-model-switcher</string>
    <string>update-models</string>
  </array>
  <key>StartInterval</key>
  <integer>86400</integer>
  <key>StandardOutPath</key>
  <string>@HOME@/Library/Logs/auto-model-switcher/update.log</string>
  <key>StandardErrorPath</key>
  <string>@HOME@/Library/Logs/auto-model-switcher/update.error.log</string>
</dict>
</plist>
